# easypbx/setup/supervisord/nextgenswitch.ini
[program:easypbx_agi]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/easypbx/artisan app:amphp
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=apache
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/easypbx/storage/logs/agi.log
stopwaitsecs=3600

[program:easypbx_queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/easypbx/artisan queue:work  --sleep=3 --tries=1 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=apache
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/easypbx/storage/logs/worker.log
stopwaitsecs=3600

[program:easypbx_schedule]
process_name=%(program_name)s_%(process_num)02d
directory=/var/www/html/easypbx
command=php artisan schedule:work
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=apache
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/easypbx/storage/logs/worker.log
stopwaitsecs=3600

[program:nextgenswitch]
process_name=%(program_name)s_%(process_num)02d
command=/usr/infosoftbd/nextgenswitch/nextgenswitchctl monitor
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=nextgenswitch
numprocs=1
redirect_stderr=true
stdout_logfile=/usr/infosoftbd/nextgenswitch/logs/switch.log
stopwaitsecs=3600

[program:nextgenswitch_firewall]
process_name=%(program_name)s_%(process_num)02d
command=/usr/infosoftbd/nextgenswitch/nextgenswitchctl firewall
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/usr/infosoftbd/nextgenswitch/logs/firewall.log
stopwaitsecs=3600